# products/wisLora/package/base-files/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
# for zl38062 reset pin
# start wisap application
msload_go.sh

exit 0
